# src/ooo/apache-ant/profile.d/ant.sh
#!/bin/sh
export ANT_HOME=/usr/lib/apache-ant
export PATH="${PATH}:${ANT_HOME}/bin"
