# build/src/apache-ant/profile.d/ant.csh
#!/bin/csh
setenv ANT_HOME /usr/lib/apache-ant
setenv PATH ${PATH}:${ANT_HOME}/bin
